# systemd/discord-bot.service
[Unit]
Description=LLM Portfolio Journal Discord Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/LLM-portfolio-project
Environment="PATH=/home/ec2-user/LLM-portfolio-project/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ec2-user/LLM-portfolio-project"
Environment="PYTHONUNBUFFERED=1"

# AWS Secrets Manager configuration
Environment="USE_AWS_SECRETS=1"
Environment="AWS_REGION=us-east-1"
Environment="AWS_SECRET_NAME=qqqAppsecrets"
Environment="AWS_RDS_SECRET_NAME=RDS/ohlcvdata"

# Use wrapper script that loads secrets before starting bot
ExecStart=/home/ec2-user/LLM-portfolio-project/.venv/bin/python scripts/start_bot_with_secrets.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/discord-bot/combined.log
StandardError=append:/var/log/discord-bot/error.log

# Resource limits
MemoryMax=500M
CPUQuota=50%

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ec2-user/LLM-portfolio-project/data
ReadWritePaths=/home/ec2-user/LLM-portfolio-project/logs
ReadWritePaths=/var/log/discord-bot

[Install]
WantedBy=multi-user.target
