# systemd/nightly-pipeline.service
[Unit]
Description=LLM Portfolio Nightly Pipeline (OHLCV, NLP, Stock Refresh)
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/llm-portfolio
Environment="PATH=/home/ubuntu/llm-portfolio/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/llm-portfolio"
Environment="PYTHONUNBUFFERED=1"

# AWS Secrets Manager configuration
Environment="USE_AWS_SECRETS=1"
Environment="AWS_REGION=us-east-1"
Environment="AWS_SECRET_NAME=llm-portfolio/production"

# Run the nightly pipeline script
# This orchestrates: OHLCV backfill, NLP batch processing, stock profile refresh
ExecStart=/home/ubuntu/llm-portfolio/.venv/bin/python scripts/nightly_pipeline.py

# Timeout for long-running job (30 minutes max)
TimeoutStartSec=1800

# Logging
StandardOutput=append:/var/log/portfolio-nightly/combined.log
StandardError=append:/var/log/portfolio-nightly/error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/llm-portfolio/data
ReadWritePaths=/home/ubuntu/llm-portfolio/logs
ReadWritePaths=/var/log/portfolio-nightly

[Install]
WantedBy=multi-user.target
