# systemd/nightly-pipeline.service
[Unit]
Description=LLM Portfolio Nightly Pipeline (OHLCV, NLP, Stock Refresh)
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/llm-portfolio
Environment="PATH=/home/ubuntu/llm-portfolio/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/llm-portfolio"
Environment="PYTHONUNBUFFERED=1"

# Matplotlib config directory (needed because ProtectHome=read-only blocks ~/.config)
Environment="MPLCONFIGDIR=/home/ubuntu/llm-portfolio/logs/matplotlib"
Environment="TMPDIR=/tmp"

# Batch output directory for NLP pipeline
Environment="BATCH_OUTPUT_DIR=/home/ubuntu/llm-portfolio/logs/batch_output"

# AWS Secrets Manager configuration (loaded from central env file)
# REQUIRED: Create /etc/llm-portfolio/llm.env first - see systemd/README.md
EnvironmentFile=/etc/llm-portfolio/llm.env

# Pre-flight check: fail fast with clear error if env file missing
ExecStartPre=/usr/bin/test -f /etc/llm-portfolio/llm.env

# Run the nightly pipeline script
# This orchestrates: OHLCV backfill, NLP batch processing, stock profile refresh
ExecStart=/home/ubuntu/llm-portfolio/.venv/bin/python scripts/nightly_pipeline.py

# Timeout for long-running job (30 minutes max)
TimeoutStartSec=1800

# Logging via journald (view with: journalctl -u nightly-pipeline.service -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llm-portfolio-nightly

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/ubuntu/llm-portfolio/data
ReadWritePaths=/home/ubuntu/llm-portfolio/logs

[Install]
WantedBy=multi-user.target
